# system-devnull/udev/99-hypr-power.rules
# Switch eDP-1 refresh (120 on AC / 60 on battery) when the adapter is plugged/unplugged.
# Runs in the user's Hyprland session via runuser; the script derives XDG_RUNTIME_DIR itself.
ACTION=="change", SUBSYSTEM=="power_supply", KERNEL=="AC0", ATTR{type}=="Mains", \
    RUN+="/usr/bin/runuser -u a1nu -- /home/a1nu/.config/hypr/scripts/hypr-refresh.sh"
