# systemd-user/.config/systemd/user/xsession.target
[Unit]
Description=X session target (startx / bspwm)
