# systemd-user/.config/systemd/user/xss-lock.service
[Unit]
Description=Screen locker on suspend (xss-lock + xsecurelock)
After=xsession.target
PartOf=xsession.target

[Service]
Type=simple
ExecStart=%h/.local/bin/xss-lock-xsecurelock
Restart=on-failure
RestartSec=1
Environment=XSECURELOCK_SAVER=/usr/lib/xsecurelock/saver_blank
Environment=XSECURELOCK_DIMMER=/usr/lib/xsecurelock/dimmer
Environment=XSECURELOCK_BLANK_TIMEOUT=5

[Install]
WantedBy=xsession.target
