# systemd-user/.config/systemd/user/xss-lock.service
[Unit]
Description=Screen locker on suspend (xss-lock + i3lock)
After=xsession.target
PartOf=xsession.target

[Service]
Type=simple
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
ExecStart=/usr/bin/xss-lock -- /usr/bin/i3lock -n -c 000000
Restart=always
RestartSec=2

[Install]
WantedBy=xsession.target

